# root/etc/profile.d/env.sh
export PATH="$PATH:/usr/local/bin/dwmscripts"
export PATH="$PATH:/usr/local/bin/dwmscripts/statusbar"
