# root/etc/profile.d/env.sh
# Use this file for general PATH exports
export TERMINAL=st

#export PATH="$PATH:/usr/local/bin/dwmscripts"
