# root/etc/profile.d/env.sh
# Use this file for general PATH exports
export TERMINAL=st
export GDK_SCALE=1
export PATH=$PATH:/usr/local/go/bin:/opt/software/node-v14.20.1-linux-x64/bin
export GOPRIVATE=github.com/liberetech
export GOPATH=$(go env GOPATH) # Or any explicit path
