# root/etc/profile.d/env.sh
# Use this file for general PATH exports
export TERMINAL=st
export GDK_SCALE=1
